# supervisord.conf
[supervisord]
nodaemon=true

[program:SSH]
command=/usr/sbin/sshd -D

[program:Cron]
command=/usr/sbin/cron -f